# templates/vsftpd.service.j2
[Unit]
Description=Vsftpd %I ftp daemon
After=network.target
PartOf=vsftpd.target

[Service]
Type=forking
ExecStart=/usr/sbin/vsftpd /etc/vsftpd/%i/vsftpd.conf

[Install]
WantedBy=multi-user.target
